# setup/systemd/buildStatistics.timer
[Unit]
Description=Run CDRTool buildStatistics every 5 mins

[Timer]
OnCalendar=*-*-* *:0/5

[Install]
WantedBy=timers.target