# setup/systemd/purgeSIPTrace.timer
[Unit]
Description=Run CDRTool purge SIP Trace each day

[Timer]
OnCalendar=*-*-* 03:20:00
RandomizedDelaySec=900
