# setup/systemd/quotaDailyReset.timer
[Unit]
Description=Run CDRTool daily quota reset each day

[Timer]
OnCalendar=*-*-* 00:10:00
RandomizedDelaySec=900

