# debian/purgeSIPTrace.service

[Unit]
Description=CDRTool purge siptraces from central database

[Service]
Type=oneshot
IgnoreSIGPIPE=false
KillMode=process
ExecStart=/var/www/CDRTool/scripts/OpenSIPS/purgeSIPTrace.php
StandardOutput=null

