# debian/purgeTables.timer
[Unit]
Description=Run CDRTool purgeTables each day

[Timer]
OnCalendar=*-*-* 04:00:00
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
