# debian/rotateTables.timer
[Unit]
Description=Run CDRTool rotate Tables each month

[Timer]
OnCalendar=*-*-1 03:00:00
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
